# helpers/99-wifi.rules
KERNEL=="wlan*", RUN+="/bin/sh -c /opt/piratebox/bin/pi-starter.sh"
